# tor-autosocket.socket
[Socket]
Accept=false
ListenStream=9050
Service=tor-autosocket.service

[Install]
WantedBy=sockets.target
